# backup.cron
0 13-15 * * * PATH="/usr/local/bin:$PATH" /Users/andrew/.backup/backup.sh 2>> /Users/andrew/.backup/log/backup.err >> /Users/andrew/.backup/log/backup.log
